# qjackctl.service
[Unit]
Description=QJackCTL

[Service]
ExecStart=/usr/bin/qjackctl

[Install]
WantedBy=user-apps.target
