# xfsettingsd.service
[Unit]
Description=XFCE Settings Daemon

[Service]
Type=forking
ExecStart=/usr/bin/xfsettingsd

[Install]
WantedBy=default.target
